# docker/wordpress/supervisor/supervisor.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:php-fpm]
command=php-fpm --nodaemonize
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startretries=3
autorestart=true

[program:tail-error-log]
command=tail -f /var/log/php-fpm/error.log
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startretries=3
autorestart=true

[program:sync-error-log]
command=lsyncd -logfile /dev/stdout -nodaemon -delay 0.5 -rsync /var/log/php-fpm /home/www-data/app/wp-content/logs/php-fpm
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
startretries=3
autorestart=true

[program:lsyncd]
command=lsyncd -logfile /dev/stdout -nodaemon -delay 0.5 -rsync /home/www-data/app /app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
autostart=true
autorestart=unexpected
startretries=3
